# roles/smartctl_exporter/templates/smartctl_exporter.socket.j2
[Unit]
Description=Smartctl Exporter

[Socket]
ListenStream=9633

[Install]
WantedBy=sockets.target